# com.hermes.weibo.downloader.800.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.hermes.weibo.downloader.800</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/henryking/Desktop/AI/Skills/MySkills/weibo-downloader/weibo_downloader_837.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key><integer>8</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <key>RunAtLoad</key><false/>
    <key>StandardOutPath</key><string>/Users/henryking/.hermes/weibo_downloader_800.out</string>
    <key>StandardErrorPath</key><string>/Users/henryking/.hermes/weibo_downloader_800.err</string>
  </dict>
</plist>
